# root/etc/services.d/plex/finish
#!/bin/sh
# Cleanup on plex exit
rm -f /tmp/plex_claim_token
